# .service
[Unit]
Description=Zocratic MMA Backend
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/zocraticmma
#Environment="PATH=/home/ubuntu/zocraticmma/venv/bin"
ExecStart=/home/ubuntu/zocraticmma/venv/bin/gunicorn backend.api.main:app --bind 0.0.0.0:8000
Restart=always

[Install]
WantedBy=multi-user.target
